# dist/bmc-vmi-ca-manager.service
[Unit]
Description=BMC VMI CA authority manager

[Service]
ExecStart=/usr/bin/bmc-vmi-ca
Restart=always
UMask=0007

Type=dbus
BusName=xyz.openbmc_project.Certs.ca.authority.Manager

[Install]
WantedBy=multi-user.target
