# dist/bmc-vmi-ca-manager.service
[Unit]
Description=BMC VMI CA authority manager

[Service]
ExecStart=/usr/bin/env bmc-vmi-ca
SyslogIdentifier=bmc-vmi-ca
Restart=always
DynamicUser=yes
Group=certificate-manager
SupplementaryGroups=phosphor-logging
UMask=0007

Type=dbus
BusName=xyz.openbmc_project.Certs.ca.authority.Manager

[Install]
WantedBy=multi-user.target
