# dist/phosphor-certificate-manager@.service
[Unit]
Description=Phosphor certificate manager for %I

[Service]
Environment=UNIT=""
EnvironmentFile=/usr/share/phosphor-certificate-manager/%I
ExecStart=/usr/bin/phosphor-certificate-manager --endpoint ${ENDPOINT} --path ${CERTPATH} --type ${TYPE} --unit ${UNIT}
Restart=always
UMask=0007

[Install]
WantedBy=multi-user.target
